# tunnel_3770_coder1.conf
description "Service Tunnel for 3770 Coder1"
author "Pin Lin"

start on started sshd
stop on stopping sshd

setuid {{user}}

script
    while true
    do
        ssh -o "ExitOnForwardFailure yes" -NR 38701:192.168.168.101:8701 ntut.com.tw
        sleep 1s
    done
end script

respawn
respawn limit 1000 1
