# tunnel_ds216j_dsm_http.conf
description "Service Tunnel for DS216J DSM HTTP"
author "Pin Lin"

start on started sshd
stop on stopping sshd

setuid {{user}}

script
    while true
    do
        ssh -o "ExitOnForwardFailure yes" -NR 5002:localhost:5000 ntut.com.tw
        sleep 1s
    done
end script

respawn
respawn limit 1000 1
