# tunnel_3770_code3.conf
description "Service Tunnel for 3770 Code3"
author "Pin Lin"

start on started sshd
stop on stopping sshd

setuid {{user}}

script
    while true
    do
        ssh -o "ExitOnForwardFailure yes" -NR 38703:192.168.168.101:8703 ntut.com.tw
        sleep 1s
    done
end script

respawn
respawn limit 1000 1
