# tunnel_ds216j_mariadb.conf
description "Service Tunnel for DS216J MariaDB"
author "Pin Lin"

start on started sshd
stop on stopping sshd

setuid {{user}}

script
    while true
    do
        ssh -o "ExitOnForwardFailure yes" -NR 3307:localhost:3306 ntut.com.tw
        sleep 1s
    done
end script

respawn
respawn limit 1000 1
